# yapo_php/supervisord-apache2.conf
[program:apache2]
command=/yapo/start-apache2.sh
numprocs=1
autostart=true
autorestart=true
