# jenkins/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
command=/bin/bash -c "java -jar /usr/share/jenkins/jenkins.war"
